# .docker-compose/web-server/etc/cron.d/100-docker-cron
* * * * * root curl -u federation:secret -X GET http://localhost/bin/federate/cleanup >> /data/poirot/cron.log 2>&1
